# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:crm]
command=python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/srv/crm
autostart=true
autorestart=true
stderr_logfile=/tmp/crm.err.log
stdout_logfile=/tmp/crm.out.log
environment=CRM_CALLBACK_SECRET="prod-shared-secret",CRM_CHANNEL_BASE_URL="http://127.0.0.1:8001",CRM_CRM_PUBLIC_URL="http://127.0.0.1:8000",CRM_BUS="memory",CRM_DATABASE_URL="sqlite+aiosqlite:///../crm.db"

[program:channel]
command=python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8001
directory=/srv/channel
autostart=true
autorestart=true
stderr_logfile=/tmp/channel.err.log
stdout_logfile=/tmp/channel.out.log
environment=CHANNEL_CALLBACK_SECRET="prod-shared-secret"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/tmp/nginx.err.log
stdout_logfile=/tmp/nginx.out.log

